# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

# Expose supervisord over a unix socket so supervisorctl (and the health check
# endpoint, running as www-data) can query program state. chown to root:www-data
# with chmod 0770 keeps the socket off-limits to anything else in the container.
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770
chown=root:www-data

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:meilisearch]
command=meilisearch --http-addr 127.0.0.1:7700 --env production --master-key %(ENV_MEILI_MASTER_KEY)s --db-path /var/lib/meilisearch/data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=5
stopsignal=TERM
stopwaitsecs=10

[program:frankenphp]
command=frankenphp run --config /etc/frankenphp/Caddyfile
user=www-data
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
priority=10
stopsignal=TERM
stopwaitsecs=30

[program:worker-messenger]
user=www-data
# --memory-limit needs a K/M/G suffix; an unsuffixed integer is parsed as
# raw bytes by Symfony's ConsumeMessagesCommand::convertToBytes(), making
# the listener fire on the worker's first idle tick.
command=php bin/console messenger:consume transactional_email notification borrow_lifecycle deck_enrichment tcgdex_sync_series tcgdex_sync_serie tcgdex_sync_set tcgdex_sync_card --time-limit=1200 --memory-limit=512M --sleep=60 --env=prod --no-debug
directory=/app
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
# Default startretries is 3, which can flip to FATAL on transient DB blips.
# Bumped because the FATAL transition now terminates the whole container.
startretries=5
priority=20
stopsignal=TERM
stopwaitsecs=30

# Watch worker-messenger and shut supervisord (PID 1) down if it reaches FATAL.
# Pairs with supervisor-fatal-listener.sh — see that script for the protocol.
# The container exiting triggers the orchestrator to spawn a fresh instance,
# which is the only way to recover from supervisord giving up on autorestart.
[eventlistener:worker-fatal-killer]
command=/usr/local/bin/supervisor-fatal-listener.sh
events=PROCESS_STATE_FATAL
autorestart=unexpected
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
